# deathclock.service
[Unit]
Description = Deathclock shutdown service
After = mediacenter.service
Wants = network-online.target

[Service]
Type = simple
ExecStart = /home/osmc/sediq/deathclock
ExecStop = /home/osmc/sediq/shutdown

[Install]
WantedBy = multi-user.target
